# deploy/bird-id-monitor.service
[Unit]
Description=Bird ID - Arlo bird-species monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=birdid
WorkingDirectory=/opt/bird-id
ExecStart=/opt/bird-id/.venv/bin/python /opt/bird-id/main.py
# Restart on crash and after clean exit (Arlo session drops, transient network).
Restart=always
RestartSec=15
# Don't get stuck in a tight crash loop.
StartLimitIntervalSec=300
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
