# deploy/bird-id-dashboard.service
[Unit]
Description=Bird ID - web dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=birdid
WorkingDirectory=/opt/bird-id
# Unbuffered so the app's print() output reaches journald in real time.
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/bird-id/.venv/bin/python /opt/bird-id/dashboard.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
